# bitnami/mariadb/testing/ol9/rootfs/opt/bitnami/scripts/cron/backup.cron
0 2 * * * /opt/bitnami/scripts/backup.sh
