# bitnami/percona-mysql/testing/8.4/ol9/rootfs/opt/bitnami/scripts/conf/mysql.ini
[program:mysql]
command=/bin/bash /opt/bitnami/scripts/mysql/entrypoint.sh /opt/bitnami/scripts/mysql/run.sh
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=unexpected
startretries=3
exitcodes=0
